# MCP-GateWay/deploy/ml-server.service
[Unit]
Description=ML Inference Server
After=network.target

[Service]
Type=simple
User=mcp
WorkingDirectory=/opt/mcp-gateway/ml-server
ExecStart=/opt/mcp-gateway/ml-server/venv/bin/python3 ml_inference_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# 환경 변수
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target

